# scripts/backup.timer
[Unit]
Description=AI Trading System Daily Backup Timer
Requires=backup.service

[Timer]
# Run daily at 2:00 AM
OnCalendar=daily
OnCalendar=02:00:00
Persistent=true

# Random delay to avoid load spikes
RandomizedDelaySec=300

[Install]
WantedBy=timers.target